# systemd/wg-netns.target
[Unit]
Description=WireGuard Tunnels in Network Namespaces
