# systemd/deluged.service.d/netns.conf
[Unit]
After=wg-netns.target
Requires=wg-netns.target

[Service]
NetworkNamespacePath=/run/netns/vpn
BindReadOnlyPaths=/etc/netns/vpn/resolv.conf:/etc/resolv.conf:norbind
